# deploy/lumi.service
[Unit]
Description=Lumi 智能家居统一层
After=network.target

[Service]
Type=simple
ExecStart=/home/amaranth/.local/bin/lumi
Restart=on-failure
RestartSec=5s
WorkingDirectory=/home/amaranth

# 环境变量可在此覆盖
# Environment=LUMI_HA_BASE_URL=http://192.168.5.184:8123
# Environment=LUMI_SERVER_TOKEN=your-token-here

# 代理设置（Lumi 内部会清掉 HA 请求的代理，这里设置的是全局默认）
Environment=NO_PROXY=127.0.0.1,localhost,192.168.0.0/16

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
